# schedule.cron
0 23 * * * python3 /home/opc/OCI-TagCompute/OCI-TagCompute.py
